# container-files/etc/supervisor.d/mariadb.conf
[program:mariadb]
directory=/var/lib/mysql/
user=mysql
command=/usr/bin/mysqld_safe --datadir='/var/lib/mysql'
autostart=true
autorestart=true
redirect_stderr=true
